# package/trousers/tcsd.service
[Unit]
Description=TCS resource manager
After=network.target

[Service]
ExecStart=/sbin/tcsd -f

[Install]
WantedBy=basic.target
